# home/.config/systemd/user/gitsign-credential-cache.service
[Unit]
Description=GitSign credential cache

[Service]
Type=simple
ExecStart=%h/.local/go/bin/gitsign-credential-cache

[Install]
WantedBy=default.target
